# .config/systemd/user/powermanager.service
[Service]
Environment=STARTUP_WAIT=30s
ExecStart=/home/huy/bin/powermanager

[Install]
WantedBy=default.target
